# build/void-include/etc/sv/native-qemu-agent/finish
#!/bin/sh
# If the agent exits cleanly with poweroff intent, halt; otherwise respawn via runit.
exit 0
